# valve.conf
[program:valve]
command = python -u /home/pi/scripts/controller.py
user = root
stopsignal = INT
autostart = true  
autorestart = true  
startsecs = 1  
stderr_logfile = /var/log/valve.err.log  ;标准错误日志
stdout_logfile = /var/log/valve.out.log  ;标准输出日志
